# layout/DEBIAN/postinst
#!/bin/sh

launchctl unload /Library/LaunchDaemons/com.torrekie.screendumpd.plist
launchctl load /Library/LaunchDaemons/com.torrekie.screendumpd.plist

exit 0;
